# resources/systemd/voicium.service
[Unit]
Description=Voicium push-to-talk dictation daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/voicium daemon
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
